# openclaw-hybrid/systemd/openclaw-dispatcher.service
[Unit]
Description=OpenClaw Queue Dispatcher
After=network-online.target openclaw-gateway.service openclaw-bridge.service
Wants=network-online.target

[Service]
Type=oneshot
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw-hybrid
EnvironmentFile=/opt/openclaw-hybrid/config/env.cloud
ExecStart=/usr/bin/env node /opt/openclaw-hybrid/scripts/dispatcher.mjs
